# etc/examples/rc.local
#!/bin/sh
#
# $OpenBSD: rc.local,v 1.3 2014/07/14 10:15:33 ajacoutot Exp $
#
# Site-specific startup actions, daemons, and other things which
# can be done AFTER your system goes into securemode.  For actions
# which should be done BEFORE your system has gone into securemode
# please see /etc/rc.securelevel.
#
