# etc/rc.d/bpflogd
#!/bin/ksh
#
# $OpenBSD: bpflogd,v 1.1 2025/05/16 06:27:20 dlg Exp $

daemon="/usr/sbin/bpflogd"

. /etc/rc.d/rc.subr

rc_cmd $1
