# etc/rc.d/ldapd
#!/bin/ksh
#
# $OpenBSD: ldapd,v 1.3 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/ldapd"

. /etc/rc.d/rc.subr

pexp="ldapd: ldap server.*"

rc_cmd $1
