# etc/rc.d/ldpd
#!/bin/ksh
#
# $OpenBSD: ldpd,v 1.3 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/ldpd"

. /etc/rc.d/rc.subr

rc_cmd $1
