# etc/rc.d/lpd
#!/bin/ksh
#
# $OpenBSD: lpd,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/lpd"

. /etc/rc.d/rc.subr

rc_cmd $1
