# etc/rc.d/rad
#!/bin/ksh
#
# $OpenBSD: rad,v 1.2 2022/10/14 10:41:43 kn Exp $

daemon="/usr/sbin/rad"

. /etc/rc.d/rc.subr

rc_configtest() {
	# use rc_exec here since daemon_flags may contain arguments with spaces
	rc_exec "${daemon} -n ${daemon_flags}"
}

rc_cmd $1
