# etc/rc.d/rbootd
#!/bin/ksh
#
# $OpenBSD: rbootd,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/rbootd"

. /etc/rc.d/rc.subr

rc_cmd $1
