# etc/rc.d/route6d
#!/bin/ksh
#
# $OpenBSD: route6d,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/route6d"

. /etc/rc.d/rc.subr

rc_cmd $1
