# etc/rc.d/sensorsd
#!/bin/ksh
#
# $OpenBSD: sensorsd,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/sensorsd"

. /etc/rc.d/rc.subr

rc_cmd $1
