# etc/rc.d/sndiod
#!/bin/ksh
#
# $OpenBSD: sndiod,v 1.6 2019/11/22 06:31:30 claudio Exp $

daemon="/usr/bin/sndiod"

. /etc/rc.d/rc.subr

rc_cmd $1
