# etc/rc.d/watchdogd
#!/bin/ksh
#
# $OpenBSD: watchdogd,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/watchdogd"

. /etc/rc.d/rc.subr

rc_cmd $1
