# debian/linux-cloud-tools-common.hv-fcopy-daemon.service
# On Azure/Hyper-V systems start the hv_fcopy_daemon
#
# author "Andy Whitcroft <apw@canonical.com>"
[Unit]
Description=Hyper-V File Copy Protocol Daemon
ConditionVirtualization=microsoft
ConditionPathExists=/dev/vmbus/hv_fcopy
BindsTo=sys-devices-virtual-misc-vmbus\x21hv_fcopy.device

[Service]
ExecStart=/usr/sbin/hv_fcopy_daemon -n

[Install]
WantedBy=multi-user.target
