# debian/linux-cloud-tools-common.hv-vss-daemon.service
# On Azure/Hyper-V systems start the hv_vss_daemon
#
# author "Andy Whitcroft <apw@canonical.com>"
[Unit]
Description=Hyper-V VSS Protocol Daemon
ConditionVirtualization=microsoft

[Service]
ExecStart=/usr/sbin/hv_vss_daemon -n

[Install]
WantedBy=multi-user.target
