# contrib/libpcap/org.tcpdump.chmod_bpf.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.tcpdump.chmod_bpf</string>
	<key>RunAtLoad</key>
	<true/>
	<key>Program</key>
	<string>/usr/local/bin/chmod_bpf</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/chmod_bpf</string>
	</array>
</dict>
</plist>
