# libexec/rc/rc.d/autounmountd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: autounmountd
# REQUIRE: FILESYSTEMS
# BEFORE: DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="autounmountd"
desc="daemon unmounting automounted filesystems"
rcvar="autofs_enable"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
