# libexec/rc/rc.d/bridge
#!/bin/sh
#
# Copyright (c) 2006 The FreeBSD Project. All rights reserved.
#
# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions
# are met:
# 1. Redistributions of source code must retain the above copyright
#    notice, this list of conditions and the following disclaimer.
# 2. Redistributions in binary form must reproduce the above copyright
#    notice, this list of conditions and the following disclaimer in the
#    documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE PROJECT ``AS IS'' AND ANY EXPRESS OR
# IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE IMPLIED WARRANTIES
# OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE ARE DISCLAIMED.
# IN NO EVENT SHALL THE PROJECT BE LIABLE FOR ANY DIRECT, INDIRECT,
# INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL DAMAGES (INCLUDING, BUT
# NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS OR SERVICES; LOSS OF USE,
# DATA, OR PROFITS; OR BUSINESS INTERRUPTION) HOWEVER CAUSED AND ON ANY
# THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT LIABILITY, OR TORT
# (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY OUT OF THE USE OF
# THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF SUCH DAMAGE.
#
# $FreeBSD$
#

# PROVIDE: bridge
# REQUIRE: netif ppp stf
# KEYWORD: nojail

. /etc/rc.subr
. /etc/network.subr

name="bridge"
desc="Network bridge setup"
start_cmd="bridge_start"
stop_cmd="bridge_stop"
cmd=""

glob_int() {
	case "$1" in
		$2 ) true ;;
		* ) false ;;
	esac
}

bridge_test() {
	bridge=$1
	iface=$2

	eval interfaces=\$autobridge_${bridge}
	if [ -n "${interfaces}" ]; then
		for i in ${interfaces}; do
			if glob_int $iface $i ; then
				ifconfig $bridge $cmd $iface > /dev/null 2>&1
				return
			fi
		done
	fi
}

autobridge()
{
	if [ -n "${autobridge_interfaces}" ]; then
		if [ -z "$iflist" ]; then
			# We're operating as a general network start routine.
			iflist="`list_net_interfaces`"
		fi

		for br in ${autobridge_interfaces}; do
			for i in $iflist; do
				bridge_test $br $i
			done
		done
	fi
}

bridge_start()
{
	cmd="addm"
	autobridge
}

bridge_stop()
{
	cmd="deletem"
	autobridge
}

iflist=$2

load_rc_config $name
run_rc_command "$1"
