# libexec/rc/rc.d/kadmind
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: kadmind
# REQUIRE: kdc
# KEYWORD: shutdown

. /etc/rc.subr

name=kadmind
desc="Server for administrative access to Kerberos database"
rcvar=${name}_enable
required_vars=kdc_enable
start_precmd=${name}_start_precmd

set_rcvar_obsolete kadmind5_server_enable kadmind_enable
set_rcvar_obsolete kadmind5_server kadmind_program
set_rcvar_obsolete kerberos5_server_enable kdc_enable

kadmind_start_precmd()
{

	command_args="$command_args &"
}

load_rc_config $name
run_rc_command "$1"
