# libexec/rc/rc.d/nfsuserd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nfsuserd
# REQUIRE: NETWORKING
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="nfsuserd"
desc="Load user and group information into the kernel for NFSv4 services and support manage-gids for all NFS versions"
rcvar="nfsuserd_enable"
command="/usr/sbin/${name}"
sig_stop="USR1"

load_rc_config $name
start_precmd="nfsuserd_precmd"

nfsuserd_precmd()
{
	if checkyesno nfs_server_managegids; then
		rc_flags="-manage-gids ${nfsuserd_flags}"
	fi
	return 0
}

run_rc_command "$1"
