# libexec/rc/rc.d/resolv
#!/bin/sh
#
# Copyright (c) 1999  Matt Dillon
# All rights reserved.
#
# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions
# are met:
# 1. Redistributions of source code must retain the above copyright
#    notice, this list of conditions and the following disclaimer.
# 2. Redistributions in binary form must reproduce the above copyright
#    notice, this list of conditions and the following disclaimer in the
#    documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE AUTHOR AND CONTRIBUTORS ``AS IS'' AND
# ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
# IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE
# ARE DISCLAIMED.  IN NO EVENT SHALL THE AUTHOR OR CONTRIBUTORS BE LIABLE
# FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL
# DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS
# OR SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS INTERRUPTION)
# HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT
# LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY
# OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF
# SUCH DAMAGE.
#
# $FreeBSD$
#

# PROVIDE: resolv
# REQUIRE: netif FILESYSTEMS
# KEYWORD: nojail

. /etc/rc.subr

name="resolv"
rcvar="resolv_enable"
desc="Create /etc/resolv.conf from kenv"
start_cmd="${name}_start"
stop_cmd=':'

# if the info is available via dhcp/kenv
# build the resolv.conf
#
resolv_start()
{
	if [ -n "`/bin/kenv dhcp.domain-name-servers 2> /dev/null`" ]; then
		interface="`/bin/kenv boot.netif.name`"
		(
		if [ -n "`/bin/kenv dhcp.domain-name 2> /dev/null`" ]; then
			echo domain `/bin/kenv dhcp.domain-name`
		fi
	
		set -- `/bin/kenv dhcp.domain-name-servers`
		for ns in `IFS=','; echo $*`; do
			echo nameserver $ns
		done
		) | /sbin/resolvconf -a ${interface}:dhcp4
	fi
}

load_rc_config $name
run_rc_command "$1"
