# libexec/rc/rc.d/securelevel
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: securelevel
# REQUIRE: adjkerntz ipfw ipfilter pf

. /etc/rc.subr

name="securelevel"
desc="Securelevel configuration"
rcvar='kern_securelevel_enable'
start_cmd="securelevel_start"
stop_cmd=":"

# Last chance to set sysctl variables that failed the first time.
#
/etc/rc.d/sysctl lastload

securelevel_start()
{
	if [ ${kern_securelevel} -ge 0 ]; then
		echo 'Raising kernel security level: '
		${SYSCTL} kern.securelevel=${kern_securelevel}
	fi
}

load_rc_config $name
run_rc_command "$1"
