# libexec/rc/rc.d/swap
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: swap
# REQUIRE: disks
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="swap"
desc="Setup swap space"
start_cmd='/sbin/swapon -aq'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
