# contrib/wpa/wpa_supplicant/systemd/wpa_supplicant.service.in
[Unit]
Description=WPA supplicant
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=@DBUS_INTERFACE@
ExecStart=@BINDIR@/wpa_supplicant -u

[Install]
WantedBy=multi-user.target
Alias=dbus-@DBUS_INTERFACE@.service
